# examples/nimage.service
[Unit]
Description=Image Service With Deep Learning
Documentation=https://githubs.com/delldu/NImage
After=syslog.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nimage --server
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=10
Restart=always
StartLimitInterval=60s
KillMode=process

[Install]
WantedBy=multi-user.target
